# solvy_gui_v3/etc/xdg/autostart/solvy-gui.desktop

[Desktop Entry]
Type=Application
Name=Solvy AI Assistant
Exec=/usr/bin/solvy-gui
X-GNOME-Autostart-enabled=true
