# solvy_daemon_v3/systemd/solvy.service
[Unit]
Description=Solvy Daemon
[Service]
ExecStart=/usr/bin/solvyd
[Install]
WantedBy=multi-user.target
