# solvy-deb/DEBIAN/postrm
#!/bin/bash
systemctl disable solvy.service || true
systemctl daemon-reload
exit 0
