# solvy/service/solvy.service
[Unit]
Description=Solvy Daemon v2 Real
